# systemd/rss-cron-articles.service
[Unit]
Description=RSS Translator — 文章正文预缓存
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/rss-translator
EnvironmentFile=/opt/rss-translator/.env
ExecStart=/usr/bin/node /opt/rss-translator/cron-vps.js articles

# 安全加固
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/rss-translator
RestrictAddressFamilies=AF_INET AF_INET6

# 日志 → journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rss-cron-articles
